# deploy/droplet/app-execution.service
[Unit]
Description=Polymarket Signals App Execution
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/polymarket-signals
EnvironmentFile=/opt/polymarket-signals/.env.app-execution
ExecStart=/usr/bin/node /opt/polymarket-signals/dist/server/server/src/app-execution/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
