# deploy/droplet/market-intelligence.service
[Unit]
Description=Polymarket Signals Market Intelligence
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/polymarket-signals
EnvironmentFile=/opt/polymarket-signals/.env.market-intelligence
ExecStart=/usr/bin/node --max-old-space-size=4096 /opt/polymarket-signals/dist/server/server/src/market-intelligence/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
