# deploy/droplet/signalite.service
[Unit]
Description=Polymarket Signals
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_DIR__/.env
ExecStart=/usr/bin/node dist/server/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

